# etc/systemd/system/ssh.service.d/override.conf
[Service]

# Nginx & sshd should basically never die.
Restart = always
OOMScoreAdjust = -800
