# etc/systemd/system/php8.4-fpm.service.d/override.conf
[Service]

# MariaDB & PHP-FPM have 2nd-highest importance after Nginx & sshd.
Restart = on-abnormal
OOMScoreAdjust = -600

MemoryMax = 20G
MemoryHigh = 16G

LimitNOFILE = 65535:524288
