# etc/systemd/system/elasticsearch.service.d/override.conf
[Service]

# Sacrifice search for memory in dire times.
Restart = no
OOMScoreAdjust = 400

MemoryMax = 16G
MemoryHigh = 12G

LimitNOFILE = 65535:524288
