# etc/systemd/system/mariadb.service.d/override.conf
[Service]

# MariaDB & PHP-FPM have 2nd-highest importance after Nginx & sshd.
Restart = on-abnormal
OOMScoreAdjust = -600

MemoryMax = 16G
MemoryHigh = 12G

LimitNOFILE = 65535:524288
