# etc/systemd/system/nginx.service.d/override.conf
[Service]

# Nginx & sshd should basically never die.
Restart = always
OOMScoreAdjust = -800

MemoryMax = 24G
MemoryHigh = 20G

LimitNOFILE = 65535:524288
